# docker/server/config/supervisord/supervisord_cron.conf
[program:cron]
command = cron -f -L 15
startsecs = 5
stopwaitsecs = 3600
stopasgroup = false
killasgroup = true